# deploy/zenohd.service
[Unit]
Description=Body Zenoh router
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=bruce
Group=bruce
WorkingDirectory=/home/bruce/zenoh/1.9.0
ExecStart=/home/bruce/zenoh/1.9.0/zenohd -c /home/bruce/Body/deploy/zenohd-router.json
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
